# greetd/greetd.service.d/lifegreet.conf
# greetd drop-in for the lifegreet greeter (installed by greeter-install.sh
# to /etc/systemd/system/greetd.service.d/lifegreet.conf).

[Unit]
# A login screen must never permanently give up. Without this, five quick
# greeter failures put greetd into start-limit-hit and the machine boots to
# a dead vt1 until someone logs in on a TTY and fixes it by hand.
StartLimitIntervalSec=0

[Service]
# lifegreet mlockalls its whole address space so password material can never
# hit swap — ~112 MiB locked at 1080p (frame buffers scale with the output).
# systemd's default LimitMEMLOCK of 8 MiB makes every mmap past the cap fail
# with EAGAIN once mlockall(MCL_FUTURE) is armed; at boot that killed the
# greeter before it drew a frame ("greeter exited without creating a
# session"). The greeter is the only thing this service runs — let it lock.
LimitMEMLOCK=infinity
